# server_src/vpn_bot/goonvpn-bot.service
[Unit]
Description=GoonVPN Telegram Bot
After=network.target

[Service]
WorkingDirectory=/root/vpn_bot
ExecStart=/root/vpn_bot/venv/bin/python3 bot.py
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
